# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err

[program:backend]
command=node /app/backend/index.js
directory=/app/backend
autostart=true
autorestart=true
startretries=10
numprocs=1
startsecs=10
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.err
environment=
    NODE_ENV="production",
    PORT=3001,
    HUGGINGFACE_API_KEY="%(ENV_HUGGINGFACE_API_KEY)s" 